# pia-wg-reload.timer
[Unit]
Description=Periodically reload the PIA WireGuard VPN Service

[Timer]
OnBootSec=3m
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target