# pia-wg-reload.service
[Unit]
Description=Reload PIA WireGuard VPN Service

[Service]
Type=oneshot
ExecStart=/opt/pia-wg/pia-wg.sh -r