# pia-wg.service
[Unit]
Description=PIA WireGuard VPN Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/pia-wg/pia-wg.sh -f
ExecReload=/opt/pia-wg/pia-wg.sh -r
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target