# pia-wg-reload.timer
[Unit]
Description=Reload PIA WireGuard VPN Service every 10 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target